# flask_app_v2.service
[Unit]
Description=flask app service - App con flask y Gunicorn
After=network.target
[Service]
User=vagrant
Group=www-data
Environment="PATH=/home/vagrant/.local/share/virtualenvs/www-elNfpBxP/bin"
WorkingDirectory=/var/www/msdocs-python-flask-webapp-quickstart
ExecStart=/home/vagrant/.local/share/virtualenvs/www-elNfpBxP/bin/gunicorn --workers 3 --bind unix:/var/www/msdocs-python-flask-webapp-quickstart/ms.sock wsgi:app

[Install]
WantedBy=multi-user.target